# systemd/cpuminerd.service
[Unit]
Description="Verium Miner Daemon"
After=network.target

[Service]
Type=forking
User=miner
ExecStart=/usr/local/bin/start-cpuminerd.sh
ExecStop=/usr/local/bin/stop-cpuminerd.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
